# config/autostart/chromium.desktop
[Desktop Entry]
Name=Chromium
Exec=chromium
Type=Application

