# config/autostart/discord.desktop
[Desktop Entry]
Name=Discord
Exec=discord
Type=Application

