# config/autostart/easyeffects-service.desktop
[Desktop Entry]
Name=EasyEffects
Exec=easyeffects --gapplication-service
Type=Application

