# config/autostart/imwheel.desktop
[Desktop Entry]
Name=IMWheel
Exec=imwheel -b 45
Type=Application

