# config/autostart/kitty.desktop
[Desktop Entry]
Name=kitty
Exec=kitty
Type=Application

