# config/autostart/motrix.desktop
[Desktop Entry]
Name=Motrix
Exec=motrix
Type=Application

