# config/autostart/nm-applet.desktop
[Desktop Entry]
Name=Network
Exec=nm-applet
Type=Application

