# config/autostart/qBittorrent.desktop
[Desktop Entry]
Name=qBittorrent
Exec=qbittorrent
Type=Application

