# config/autostart/steam.desktop
[Desktop Entry]
Name=Steam
Exec=steam-runtime
Type=Application

