# config/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Exec=sxhkd
Type=Application

